# systemdupdater/systemdupdater.service
[Unit]
Description=Systemd-updater
WantedBy=multi-user.target

[Service]
ExecStart=/usr/bin/systemdupdater
